# roles/convidat/files/convidat.path
[Unit]
Description=Deletes /var/lib/AccountsService/icons/convidat

[Path]
PathExists=/var/lib/AccountsService/icons/convidat

[Install]
WantedBy=paths.target

