# roles/report-idalu/files/report-idalu.service
[Unit]
Description=Reporta actividad del equipo utilizando IDALU.

[Service]
ExecStart=/opt/report-idalu
