# roles/report-idalu/files/report-idalu.timer
[Unit]
Description=Reporta uso del portátil indicando IDALU.

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
