# roles/vasilisa/files/vasilisa.service
[Unit]
Description=Vasilisa la sabia cuida de los equipos al arrancar y cada día.

[Service]
ExecStart=/opt/vasilisa

