# roles/vasilisa/files/vasilisa.timer
[Unit]
Description=Ejecuta Vasilisa en el arranque y cada día.

[Timer]
OnBootSec=60
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
