# roles/apagado-2130/files/notifica-apagado-2130.timer
[Unit]
Description=Notifica al usuario el apagado a las 21:30

[Timer]
OnCalendar=21:25
OnCalendar=21:29

[Install]
WantedBy=timers.target
