# tools/sirin/systemd/sirin-ansible.timer
[Unit]
Description=Procesa con ansible-playbook las peticiones pendientes

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
